# upstart_job.conf
description "Golem Queueing Infrastructure"
author "Jeremy Archer <open-source@fatlotus.com>"

start on network up
respawn

script
  python -m golem
end script